# rootfs/etc/s6-overlay/s6-rc.d/monitor-hub/run
#!/command/execlineb

python3 /monitor-hub/monitor-hub.py
```
