# rootfs/etc/s6-overlay/s6-rc.d/monitor-hub/run
#!/command/execlineb

/etc/s6-overlay/scripts/monitor-hub
